# backup/backup.cron
# import cron with $crontab -u <user> [file]

0 4 * * 1 bash -c "/home/jenkins/backup/backup_jenkins.sh >> /home/jenkins/backup/logs/backup_jenkins_$(date +%d-%m-%y).log 2>&1 && /home/jenkins/backup/rSyncBackups.sh -u jenkins -r snf-11813.ok-kno.grnetcloud.net -k /home/jenkins/.ssh/id_rsa -f /home/jenkins/backup/thinBackups -t /home/backups >> /home/jenkins/backup/logs/backup_jenkins_$(date +%d-%m-%y).log 2>&1"
